# nvidia-410-udev-workaround/restart-udevd.service
[Unit]
Description=Workaround for nvidia-410 bug

[Service]
ExecStart=/etc/systemd/system/restart-udevd.sh

[Install]
WantedBy=multi-user.target
